# cloud/etc/supervisor.conf
[supervisord]
nodaemon=true

[program:cloud_update]
command=/usr/bin/python /data/wb_cloud/wb_cloud/cloud_update.py

[program:cloud_service]
command=/usr/bin/python /data/wb_cloud/wb_cloud/server.py
